# overlays/bigscreen/etc/profile.d/enforce_plasma_style.sh
export QT_QUICK_CONTROLS_STYLE=org.kde.desktop
#export QT_QUICK_CONTROLS_MOBILE=1

