# overlays/mark2/etc/profile.d/enforce-qtvirtualkeyboard.sh
# Only set for X11
text "$XDG_SESSION_TYPE" = "x11" || return

export QT_IM_MODULE=qtvirtualkeyboard
